# debian/etc/systemd/system/ns.service
[Unit]
Description=NS Server Configuration Tool
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/usr/share/ns
ExecStart=/usr/bin/ns service enable app
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/log /tmp /var/www/ns

[Install]
WantedBy=multi-user.target