# build/debian/DEBIAN/postinst
#!/bin/bash
set -e

echo "Installing NS package..."

# Create SSL certificates if they don't exist
if [ ! -f "/etc/nginx/ssl/selfsigned.crt" ] || [ ! -f "/etc/nginx/ssl/selfsigned.key" ]; then
    echo "Generating self-signed SSL certificates..."
    mkdir -p /etc/nginx/ssl
    openssl req -x509 -nodes -days 365 -newkey rsa:2048 \
        -keyout /etc/nginx/ssl/selfsigned.key \
        -out /etc/nginx/ssl/selfsigned.crt \
        -subj "/C=US/ST=State/L=City/O=Organization/CN=localhost"
    chmod 600 /etc/nginx/ssl/selfsigned.key
    chmod 644 /etc/nginx/ssl/selfsigned.crt
fi

# Deploy Flutter web files
echo "Deploying web files..."
mkdir -p /var/www/ns
cp -r /usr/share/ns/web/* /var/www/ns/
chown -R www-data:www-data /var/www/ns
chmod -R 755 /var/www/ns




# Backup existing nginx config
if [ -f "/etc/nginx/nginx.conf" ]; then
    echo "Backing up existing nginx configuration..."
    cp /etc/nginx/nginx.conf "/etc/nginx/nginx.conf.backup.$(date +%Y%m%d_%H%M%S)"
fi

# Test nginx configuration
echo "Testing nginx configuration..."
if nginx -t; then
    echo "Nginx configuration is valid"
    systemctl reload nginx 2>/dev/null || echo "Nginx not running, will start with system"
else
    echo "ERROR: Nginx configuration is invalid!"
    echo "Restoring backup configuration..."
    if [ -f "/etc/nginx/nginx.conf.backup.*" ]; then
        cp $(ls -t /etc/nginx/nginx.conf.backup.* | head -1) /etc/nginx/nginx.conf
    fi
    exit 1
fi

# Enable and start services
systemctl daemon-reload

if [ -f "/etc/systemd/system/ns.service" ]; then
    systemctl enable ns.service
    systemctl start ns.service
    echo "NS service enabled and started"
fi

# Start nginx if not running
if ! systemctl is-active --quiet nginx; then
    systemctl start nginx
    echo "Nginx started"
fi

echo "NS installation completed successfully"
echo "Web interface available at: https://$(hostname -I | awk '{print $1}')"

exit 0